# pandora_agent/supervisord.conf
[supervisord]
nodaemon=true

[program:pandora_agent]
command=/usr/bin/pandora_agent /etc/pandora

[program:webapp]
command=/usr/local/bin/node /opt/webapp/app.js
